# smeserver-print-pdf/root/etc/e-smith/templates/etc/crontab/printpdf
*/5 * * * * printpdf /home/e-smith/files/users/printpdf/autoprint.sh
